# supervisord.conf
[supervisord]
nodaemon=true


[program:rabbit-version]
environment=PATH="/root/.virtualenvs/rabbitpy/bin"
command=python /usr/local/d8o/rabbitpy/rpc/version.py
