# supervisord.conf
[supervisord]
nodaemon=true


[program:rabbit-version]
environment=PATH="/root/.virtualenvs/rabbitpy/bin/python"
command=/bin/bash
