# supervisor/raygun_rabbit.conf
[program:raygun]
command=/bin/bash -c "PYTHONPATH=$PYTHONPATH:/usr/local/rabbitpy:/usr/local/domainsage \
                        python3.4 /usr/local/rabbitpy/service/raygun.py"
priority=150
redirect_stderr=true
stdout_logfile=syslog
